# claude-bridge.service
[Unit]
Description=ZeroRelay Claude Code bridge (claude_ai)
After=network-online.target zerorelay.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/zerorelay/claude-bridge.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
