# zerorelay.service
[Unit]
Description=ZeroRelay WebSocket broker (Tailscale only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/zerorelay/zerorelay.py --host 100.127.106.41 --port 8765
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
